# systemd/user/flexget-anime.timer
[Unit]
Description=Run Flexget RSS Download Targets

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
