# systemd/user/udiskie.service
[Unit]
Description=Handle automounting of usb devices

[Service]
Type=simple
ExecStart=udiskie

[Install]
WantedBy=default.target
